# luci-app-ssr-lite/root/etc/uci-defaults/40_luci-ssr-lite
#!/bin/sh

mkdir -p /etc/shadowsocksr
/etc/init.d/rpcd reload
exit 0
